# roles/system/files/aide/aide.timer
[Unit]
Description=Run AIDE daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target